# vibebot-updater.service
[Unit]
Description=Check for vibebot container updates
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
Environment=IMAGE=ghcr.io/rdrake/vibebot-v8:latest
ExecStart=/bin/bash -c '\
  LOCAL=$(docker inspect --format="{{.Image}}" vibebot 2>/dev/null || echo "none"); \
  docker pull -q ${IMAGE}; \
  REMOTE=$(docker inspect --format="{{.Id}}" ${IMAGE}); \
  if [ "$LOCAL" != "$REMOTE" ]; then \
    echo "Update found: restarting vibebot..."; \
    systemctl --user restart vibebot.service; \
  else \
    echo "No update available"; \
  fi'
