# vibebot-updater.timer
[Unit]
Description=Check for vibebot updates every 15 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
